# finpulse.conf
description "Provider website"
start on runlevel [234]
stop on runlevel [06]
respawn
respawn limit 10 5

setuid nishant
setgid nishant
exec /home/nishant/public_html/provider/gunicorn_start.sh
