# finpulse.conf
description "Finpulse.in"
start on runlevel [234]
stop on runlevel [06]
respawn
respawn limit 10 5

setuid root
setgid root
exec /root/public_html/finpulse.com/gunicorn_start.sh

